# deploy/mibvendor-health.timer
[Unit]
Description=Check mibvendor production health every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=mibvendor-health.service

[Install]
WantedBy=timers.target
